# lesson_08/include/etc/systemd/system/generator.service
[Unit]
Description=Log generator for watcher

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/generator
ExecStart=/opt/generator.sh $WORD $LOG
