# lesson_08/include/etc/systemd/system/watcher.service
[Unit]
Description=My watcher service

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/watcher
ExecStart=/opt/watcher.sh $WORD $LOG
