# lesson_08/include/etc/systemd/system/watcher.service
[Unit]
Description=My watcher service

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/generator
ExecStart=/opt/watcher.sh $WORD $LOG
